# future-news.service
[Unit]
Description=Grand Lodge of Future Sight - Future News Service
After=network.target

[Service]
User=masonic
Group=masonic
WorkingDirectory=/opt/future-news
ExecStart=/usr/bin/python -m poetry run uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
SyslogIdentifier=future-news

# Security settings
ProtectSystem=full
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
